# machines/LAPTOPS/HP-i5-ProBook-6460b/core/common/services/bluetooth.service
[Unit]
Description=Bluetooth service
After=network.target

[Service]
ExecStart=/run/current-system/sw/bin/bluetoothd
Restart=always

[Install]
WantedBy=multi-user.target

